# network-hardening-check.service
[Unit]
After=nftables.service

[Service]
Type=oneshot
ExecStart=/usr/bin/network-hardening-check
